# USYVL-Office-System/files/etc/rc.local
#!/bin/sh
/usr/local/sbin/firewallutil load
